# executable/crontab
@reboot /home/pi/sdrstart.sh >sdrlog.log
