# install-tools/crontab
@reboot sudo /home/pi/sdrstart.sh > /dev/null 2>&1
