# install/crontab
@reboot sudo /home/pi/sdrstart.sh >sdrberry.log
